# runit/throttled/run
#!/bin/sh
# PYTHONUNBUFFERED=1 exec /usr/lib/throttled/venv/bin/python3 /usr/lib/throttled/throttled.py
PYTHONUNBUFFERED=1 exec /bin/python3 /usr/lib/throttled/throttled.py
